# monolith/api.ini
[program:api]
command=/usr/local/sbin/ipsec_api
redirect_stderr=true
stdout_logfile=/tmp/api.log
#stdout_logfile_maxbytes = 0
#stderr_logfile_maxbytes = 0
