# monolith/bird.ini
[program:bird]
command=bird -c /etc/bird.conf -d -s /etc/bird.d/bird.ctl
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0