# monolith/ipsec.ini
[program:strongswan]
command=ipsec start --nofork
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0