# docker/frr.ini
[program:frr]
command=/usr/local/sbin/frr.sh
redirect_stderr=true
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes = 0
#stderr_logfile_maxbytes = 0