# docker/frr.ini
[program:frr]
command=/usr/local/sbin/frr.sh
redirect_stderr=true
stdout_logfile=/opt/logs/frr.log
stdout_logfile_maxbytes = 10485760 # 10 mb
stdout_logfile_backups = 10
