# docker/haproxy.ini
[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg
redirect_stderr=true
stdout_logfile=/opt/logs/haproxy.log
#stdout_logfile_maxbytes = 0
#stderr_logfile_maxbytes = 0
