# environment/supervisor/content/site1/ipsec-01.ini
[program:ipsec-01]
command=/usr/local/sbin/ipsec.sh
environment=PHYS_IF="eth0", VRF="100", XFRM_IP="10.0.0.1 10.0.0.5", XFRM_PEER="10.0.0.2 10.0.0.6", NAT="No Yes"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0