# environment/supervisor/content/site1/ipsec-02.ini
[program:ipsec-02]
command=/usr/local/sbin/ipsec.sh
environment=PHYS_IF="eth0", VRF="101", XFRM_IP="10.1.0.1", XFRM_PEER="10.1.0.2", NAT="Yes"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0